# cron.daily
#!/bin/bash
set -e

BASE=/var/archive-kit/britney2
# some tools use relative paths, so we *need* to cd
cd $BASE

./tg-britney pkglists
./tg-britney urgencies
./tg-britney run
./tg-britney save
./tg-britney stats
